# debian/wlanpi-chat-bot.service
[Unit]
Description=Telegram chatbot the WLAN Pi

[Service]
ExecStart=/usr/sbin/chatbot

[Install]
WantedBy=multi-user.target
