# scripts/wlanpi-chat-bot.service
[Unit]
Description=WLAN Pi Chat Bot

[Service]
ExecStart=/usr/bin/python3 /opt/wlanpi-chat-bot/wlanpi-chat-bot

[Install]
WantedBy=multi-user.target